# service/configs/include/workers.ini
[program:export_worker]
command=php /data/www/dachuwang/service/task_workers/php/a.php
priortity=1
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisord/export_worker.log ; (main log file;default $CWD/supervisord.log)
stdout_logfile_maxbytes=1MB
